# deploy/facsim.timer
# facsim — daily rebuild schedule.
# Pairs with facsim.service. Persistent=true catches up a missed run
# (e.g. the box was off at the scheduled time).
[Unit]
Description=Daily rebuild of the LLM cache

[Timer]
OnCalendar=*-*-* 04:30:00
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
